# src/etc/clam/clamdscan.cron
0 1 * * 0 root /usr/bin/clamdscan --multiscan --fdpass /
